# configs/systemd/nextepc-udmd.service.in
[Unit]
Description=NextEPC UDM Daemon
After=network-online.target

[Service]
Type=simple

User=nextepc
Group=nextepc

Restart=always
ExecStart=@bindir@/nextepc-udmd -c @sysconfdir@/nextepc/udm.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
