# configs/systemd/nextepc-bsfd.service.in
[Unit]
Description=NextEPC BSF Daemon
Wants=mongodb.service mongod.service

[Service]
Type=simple

User=nextepc
Group=nextepc

Restart=always
ExecStart=@bindir@/nextepc-bsfd -c @sysconfdir@/nextepc/bsf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
